# rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
python3 /nextcloud_file_watcher/briefkasten_watcher.py
