# nextcloud_file_watcher/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
echo "Starte Nextcloud File Watcher..."
cd /app || exit 1
exec python3 briefkasten_watcher.py
